# backend/deployment/kaya-sync.timer
# Kaya Sync - systemd Timer

[Unit]
Description=Daily Kaya Climbing Data Sync Timer
Requires=kaya-sync.service

[Timer]
# Run daily at 2:00 AM
OnCalendar=daily
OnCalendar=*-*-* 02:00:00

# If system was off, run 15 minutes after boot
OnBootSec=15min

# Run 5 minutes after timer activation if missed
Persistent=true

[Install]
WantedBy=timers.target
