# backend/deployment/kaya-sync.service
# Kaya Sync - systemd Service
# This service runs the Kaya data sync job using the compiled binary

[Unit]
Description=Kaya Climbing Data Sync Service
After=network.target postgresql.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/woulder

# Load environment variables (same .env as the API)
EnvironmentFile=/opt/woulder/.env

# Run the compiled sync binary
ExecStart=/opt/woulder/kaya-sync-job --incremental --delay 3

# Logging
StandardOutput=append:/var/log/kaya-sync/sync.log
StandardError=append:/var/log/kaya-sync/error.log

# Restart on failure
Restart=on-failure
RestartSec=300

[Install]
WantedBy=multi-user.target
